# service/log/run
#!/bin/sh
exec 2>&1
exec multilog t s250000 n4 /var/log/dbus-huaweisun2000
